# etc/cron.daily/cp_ssh
#!/bin/bash

/usr/bin/rsync -av --delete ~/.ssh /mnt/c/Users/kira
