# etc/cron.daily/update_vscode_extensions_list
#!/bin/bash

DOTFILES_DIR=~/ghq/github.com/Ras96/dotfiles

code --list-extensions > $DOTFILES_DIR/.vscode-extensions.txt
